# packaging/systemd/md-preview.service
[Unit]
Description=md-preview — always-on local Markdown preview daemon
Documentation=https://github.com/hibiki-automatic/md-preview
After=network.target

[Service]
# %h expands to the user's home directory (user-unit context).
# The binary is expected at ~/.local/bin/md-preview (the conventional XDG
# location for user-installed binaries; `cargo install` or the Makefile target
# puts it there).  Adjust the path if you installed elsewhere.
ExecStart=%h/.local/bin/md-preview --daemon
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
